# templates/minecraft-server.service.conf.erb
description "Basic Minecraft server startup script"
author "David Winter <i@djw.me>"

start on startup
stop on shutdown

respawn

expect fork

exec java -Xms<%= max_ram %>M -Xmx<%= max_ram %>M -jar <%= full_path %>
